# templates/app-supervisor.conf.j2
[program:{{app_name}}]
command={{app_command}}
autostart={{app_autostart}}
autorestart={{app_autorestart}}
user={{app_user}}
directory={{app_directory}}
