# templates/app-supervisor.conf.j2
[program:{{app_name}}]
command={{app_command}}
autostart=true
autorestart=true
